# deployment/launchd/com.cryptolabinc.runeconsole.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cryptolabinc.runeconsole</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/runeconsole</string>
    <string>daemon</string>
    <string>start</string>
    <string>--config</string>
    <string>/opt/runeconsole/configs/runeconsole.conf</string>
  </array>

  <key>UserName</key>
  <string>runeconsole</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/opt/runeconsole/logs/runeconsole.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/opt/runeconsole/logs/runeconsole.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
